# crontab
@hourly		cd /data && /data/bin/python3 /data/oto_webpage.py